# app/pull_trends.cron
0 */4 * * * cd /app/; /usr/local/bin/python -m tasks.pull_trends >> /var/log/pull_trends.log 2>&1
