# app/pull_trends.cron
* * * * * root python /app/tasks/pull_trends.py >> /var/log/pull_trends.log 2>&1
